# protected/extensions/crontab/crontabs/my_crontab
25 19 18 1 3 /var/www/edusec/edusec_new_development/protected/yiic feesschedulesms 2 6 9 0 fees "dsfa dfsa dasf " 64 2 1 1
25 19 18 1 3 /var/www/edusec/edusec_new_development/protected/yiic feesschedulesms 2 6 10 40 fees "afd asd ffasd" 64 2 2 1
25 19 18 1 3 /var/www/edusec/edusec_new_development/protected/yiic attendance 2 6 10 40 attendance 64 2 1
13 12 12 11 3 /var/www/edusec/edusec_new_development/protected/yiic attendance 2 6 10 40 attendance 64 2 1
25 19 18 1 3 /var/www/edusec/edusec_new_development/protected/yiic attendance 2 6 10 40 attendance 64 2 1
25 19 18 1 3 /var/www/edusec/edusec_new_development/protected/yiic attendance 2 6 10 40 attendance 64 2 1
25 19 18 1 3 /var/www/edusec/edusec_new_development/protected/yiic attendance 2 6 10 40 attendance 64 2 1
25 19 18 1 3 /var/www/edusec/edusec_new_development/protected/yiic attendance 2 6 10 0 attendance 64 2 1
